# Cliente-Servidor/Tarea 6 tomcat/rc.local
#!/bin/bash
iptables -t nat -A OUTPUT -o lo -p tcp --dport 80 -j REDIRECT --to-port 8080
runuser -l ubuntu -c 'export JAVA_HOME=/usr/lib/jvm/java-8-openjdk-amd64//usr/lib/jvm/java-8-openjdk-amd64;export CATALINA_HOME=/home/Alejandro/apache-tomcat-8.5.84;sh $CATALINA_HOME/bin/catalina.sh start'
exit 0